# deb/before-remove.sh
#!/bin/bash
/etc/init.d/haproxy-exporter stop
update-rc.d -f haproxy-exporter remove
